# backup/crontab
# MySQL Monthly Backup Schedule
# Runs at 2:00 AM on the 1st day of every month
0 2 1 * * root /usr/local/bin/mysql-backup.sh >> /var/log/mysql-backup.log 2>&1

# Empty line required at end
